# config/kapa-dashboard.service
# =============================================================================
# systemd-Unit für die VMware Kapazitätsplanung
# =============================================================================
# Installation:
#   sudo useradd --system --home /opt/kapa --shell /usr/sbin/nologin kapa
#   sudo mkdir -p /opt/kapa && sudo cp aria_kapa.py /opt/kapa/
#   sudo chown kapa:kapa /opt/kapa/aria_kapa.py
#   # Laufzeitdaten BEWUSST außerhalb von /opt/kapa (CI/CD darf sie nicht löschen):
#   sudo mkdir -p /var/lib/kapa && sudo chown kapa:kapa /var/lib/kapa && sudo chmod 750 /var/lib/kapa
#
#   # EINE Konfigurationsdatei (nicht-geheim), Mode 640:
#   sudo mkdir -p /etc/kapa
#   sudo cp config/kapa.ini.example /etc/kapa/kapa.ini   # und Werte eintragen!
#   sudo chown root:kapa /etc/kapa/kapa.ini && sudo chmod 640 /etc/kapa/kapa.ini
#
#   # Geheimnisse als eigene Dateien (root:kapa, 0640) – Pfade stehen in der INI:
#   sudo sh -c 'echo "DAS-ARIA-PASSWORT" > /etc/kapa/aria.pass'
#   sudo chown root:kapa /etc/kapa/aria.pass && sudo chmod 640 /etc/kapa/aria.pass
#   #  (analog ad_bind.pass / smtp.pass / backup.pass, falls genutzt)
#
#   sudo cp config/kapa-dashboard.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kapa-dashboard
# Logs: journalctl -u kapa-dashboard -f
#
# Hinweis: Es gibt bewusst KEINE kapa.env mehr. Alle nicht-geheimen Werte stehen
# in kapa.ini, Passwörter in eigenen .pass-Dateien (deren Pfad die INI nennt).
# =============================================================================

[Unit]
Description=VMware Kapazitätsplanung (Aria Operations Dashboard)
Documentation=https://github.com/marcobockelbrink/vmware
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kapa
Group=kapa
WorkingDirectory=/opt/kapa

ExecStart=/usr/bin/python3 /opt/kapa/aria_kapa.py --config /etc/kapa/kapa.ini --serve

Restart=on-failure
RestartSec=10

# Härtung: nur /var/lib/kapa (die Laufzeitdaten) ist beschreibbar, sonst read-only.
# /etc/kapa (INI + .pass) bleibt lesbar (ProtectSystem=strict macht es nur read-only).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kapa
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
