# config/kapa-dashboard.en.service
# =============================================================================
# systemd unit for VMware Capacity Planning
# (English twin of kapa-dashboard.service — the unit itself is identical,
#  install either one as kapa-dashboard.service.)
# =============================================================================
# Installation:
#   sudo useradd --system --home /opt/kapa --shell /usr/sbin/nologin kapa
#   sudo mkdir -p /opt/kapa && sudo cp aria_kapa.py /opt/kapa/
#   sudo chown kapa:kapa /opt/kapa/aria_kapa.py
#   # Runtime data DELIBERATELY outside /opt/kapa (CI/CD must not delete it):
#   sudo mkdir -p /var/lib/kapa && sudo chown kapa:kapa /var/lib/kapa && sudo chmod 750 /var/lib/kapa
#
#   # ONE configuration file (non-secret), mode 640:
#   sudo mkdir -p /etc/kapa
#   sudo cp config/kapa.ini.en.example /etc/kapa/kapa.ini   # and fill in values!
#   sudo chown root:kapa /etc/kapa/kapa.ini && sudo chmod 640 /etc/kapa/kapa.ini
#
#   # Secrets as separate files (root:kapa, 0640) – paths are set in the INI:
#   sudo sh -c 'echo "THE-ARIA-PASSWORD" > /etc/kapa/aria.pass'
#   sudo chown root:kapa /etc/kapa/aria.pass && sudo chmod 640 /etc/kapa/aria.pass
#   #  (same for ad_bind.pass / smtp.pass / backup.pass, if used)
#
#   sudo cp config/kapa-dashboard.en.service /etc/systemd/system/kapa-dashboard.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kapa-dashboard
# Logs: journalctl -u kapa-dashboard -f
#
# Note: there is deliberately NO kapa.env anymore. All non-secret values live
# in kapa.ini, passwords in separate .pass files (whose paths the INI names).
# =============================================================================

[Unit]
Description=VMware Capacity Planning (Aria Operations dashboard)
Documentation=https://github.com/marcobockelbrink/vmware
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kapa
Group=kapa
WorkingDirectory=/opt/kapa

ExecStart=/usr/bin/python3 /opt/kapa/aria_kapa.py --config /etc/kapa/kapa.ini --serve

Restart=on-failure
RestartSec=10

# Hardening: only /var/lib/kapa (the runtime data) is writable, rest read-only.
# /etc/kapa (INI + .pass) stays readable (ProtectSystem=strict makes it read-only).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kapa
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
